# clockd-server.openrc
#!/sbin/openrc-run
depend()
{
	use dns
	provide clock
}
start()
{
	clockd -f -s
}
