# clockd-client.openrc
#!/sbin/openrc-run
SERVER="192.168.2.3"
depend()
{
	provide clock
}
start()
{
	clockd -S -c "$SERVER"
}
