# clockd-client.openrc
#!/sbin/openrc-run
SERVER="$(cat /etc/clockd.txt)"
depend()
{
	provide clock
}
start()
{
	clockd -S -c "$SERVER"
}
