# clockd-server.service
[Unit]
Description=Time Synchronization Server
[Service]
Type=simple
ExecStart=/usr/local/bin/clockd -s
[Install]
WantedBy=time-sync.target
