# clockd-client.service
[Unit]
Description=Time Synchronization Client
[Service]
Type=oneshot
ExecStart=/usr/local/bin/clockc.sh
[Install]
WantedBy=time-sync.target
